# lider-server-deb/debian/postinst
#!/bin/sh

set -e

if [ -d /usr/share/lider-server ]; then
	update-rc.d lider defaults
	systemctl start lider.service
fi
